# crontab/InternetMonitor
MAILTO=""
*/1 * * * * izz /home/izz/scripts/internetMonitor/mon.sh
