# crontab/InternetMonitor
MAILTO=""
*/1 * * * * svc_intmon /opt/svc_intmon/internetMonitor/mon.sh
